# scripts/com.gongqing.hermes.daily-redeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gongqing.hermes.daily-redeploy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/gongqing/projects/hermes-agent/scripts/hermes_daily_redeploy.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/gongqing/projects/hermes-agent</string>
  <key>StandardOutPath</key>
  <string>/Users/gongqing/.hermes/logs/daily-redeploy.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/gongqing/.hermes/logs/daily-redeploy.error.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
